# systemd_units/minecraft_backup.service.template
[Unit]
Description=Minecraft Server

[Service]
Type=simple
ExecStart=SCRIPTDIR/make_backup.sh --regular
