# systemd_units/minecraft_prunebackups.service.template
[Unit]
Description=Cleans Minecraft Backups so they don't get out of control

[Service]
Type=simple
ExecStart=SCRIPTDIR/prune_backup.sh
